# etc/update-motd.d/2-fortune
#!/bin/sh

cat << EOF
	      Wise words of the day:
		`/usr/games/fortune -as | sed "s/^/         /"`
EOF
